# usr/lib/systemd/system/snapborg-backup-daily@.timer
[Unit]
Description=Daily runs of snapborg backup (and prune) for the '%I' snapper config

[Timer]
OnCalendar=daily
Persistent=true
Unit=snapborg-backup@%i.service

[Install]
WantedBy=timers.target
