# etc/systemd/system/snapborg-backup@.timer
[Unit]
Description=Scheduled runs of snapborg backup for the %I snapper config

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
