# etc/systemd/system/snapborg-prune@.timer
[Unit]
Description=Scheduled runs of snapborg prune for the %I snapper config

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
