# etc/systemd/system/snapborg-prune-all.timer
[Unit]
Description=Scheduled runs of snapborg prune

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
