# etc/systemd/system/snapborg-backup-all.timer
[Unit]
Description=Scheduled runs of snapborg backup for all configs

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
